# ha_claude_agent_addon/rootfs/etc/s6-overlay/s6-rc.d/ha-claude-agent/finish
#!/command/with-contenv bashio

bashio::log.warning "HA Claude Agent server stopped. Restarting in 5s..."
sleep 5
